# xinetd.d/bitlbee
service ircd
{
        socket_type     = stream
        wait            = no
        user            = nobody
        server          = /usr/sbin/bitlbee
}
